# StarryNight/etc/systemd/system/dhclient.service
[Unit]
Description=Run dhclient for enp6s0

[Service]
ExecStart=/bin/dhclient -d enp6s0

[Install]
WantedBy=multi-user.target
