# StarryNight/etc/systemd/system/dhclient.service
[Unit]
Description=Run dhclient for Lan

[Service]
ExecStart=/bin/dhclient -d Lan

[Install]
WantedBy=multi-user.target
